# api/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=npm start
redirect_stdout=true
stdout_logfile=/var/log/backend/backend.log
stderr_logfile=/var/log/backend/backend.err.log
