# api/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=npm start
stdout_logfile=/var/log/backend
stderr_logfile=/var/log/backend
